# debian/movai-studio.prerm
#!/bin/bash
# pre remove actions
# clean unziped files

SERVICENAME=movai-studio
APP_PATH=/etc/$SERVICENAME/launch

systemctl stop $SERVICENAME
systemctl disable $SERVICENAME
rm /etc/systemd/system/$SERVICENAME.service
rm -r $APP_PATH

systemctl daemon-reload
systemctl reset-failed