# debian/movai-studio.prerm
#!/bin/bash
# pre remove actions
# clean unziped files

SERVICENAME=movai-studio
APP_PATH=/usr/share/$SERVICENAME/launch

systemctl stop $SERVICENAME
systemctl disable $SERVICENAME
rm /usr/lib/systemd/system/$SERVICENAME.service

systemctl daemon-reload
systemctl reset-failed

# Remove the package (but not its configuration files)
if [ "$1" = "remove" ]; then
    rm -r $APP_PATH
fi